# .config/autostart/start.sh.desktop
[Desktop Entry]
Exec=/home/dawn/CustomScripts/autostart/start.sh
Icon=dialog-scripts
Name=start.sh
Path=
Type=Application
X-KDE-AutostartScript=true
